# private/docker/supervisor-app.conf
[program:node]
command=node "/usr/src/app/bundle/main.js"
